# cam/mediamtx_restarter.service
[Unit]
Wants=network.target
[Service]
ExecStart=/bin/bash /home/pi/restart_mediamtx.sh
[Install]
WantedBy=multi-user.target
